# rootfs/etc/services.d/firefox/run
#!/usr/bin/with-contenv sh

firefox --headless --remote-debugging-port 9222